# vts-backend/etc/profile.d/vts.sh
# Inject directory containing VTS binaries into PATH.
export PATH=$PATH:/opt/vts/bin
